# lib/systemd/system/f5-endpoint-agent.service
[Unit]
Description=F5 endpoint BIG-IP Agent
After=syslog.target network.target
Requires=network.target

[Service]
User=neutron
ExecStart=/usr/bin/f5-endpoint-agent --log-file /var/log/neutron/f5-endpoint-agent.log --config-file /etc/neutron/neutron.conf --config-file /etc/neutron/services/f5/f5-endpoint-agent.ini
Restart=always

[Install]
WantedBy=multi-user.target